# systemd/pcts.socket
[Unit]
Description=Puppet Change Testing System

[Socket]
ListenStream=80

[Install]
WantedBy=default.target
